# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:llmopt]
command=/app/start-saas.sh
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:lastsaas]
command=/app/lastsaas-server
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=SERVER_PORT="8091",LASTSAAS_CONFIG_DIR="/app/lastsaas-config",LASTSAAS_ENV="prod",FRONTEND_STATIC_DIR="/dev/null",DATABASE_NAME="llmopt",FRONTEND_URL="https://llmopt.metavert.io",APP_NAME="LLM Optimizer"
